# arch-daemons.target
[Unit]
Description=Legacy Arch Daemons

[Install]
WantedBy=multi-user.target
